# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start

# [program:client]
# startsecs=0
# command=npm run start:docker --prefix ./client

